# deploy/paper-mes-monitor.service.example
[Unit]
Description=Paper MES health and backup monitor
After=network-online.target paper-mes.service
Wants=network-online.target paper-mes.service

[Service]
Type=oneshot
User=paper-mes
Group=paper-mes
Environment=MONITOR_ENV_FILE=/etc/paper-mes/monitor.env
ExecStart=/usr/local/bin/monitor-paper-mes
StateDirectory=paper-mes
StateDirectoryMode=0750
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
